# ops/systemd/mspmetro-backend.service
[Unit]
Description=MSPMetro backend (Flask)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=caddy
Group=caddy
WorkingDirectory=/opt/mspmetro-ssr/cityfeed_pull
EnvironmentFile=-/etc/default/mspmetro-backend
ExecStart=/opt/mspmetro-ssr/cityfeed_pull/backend/.venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 2 mspmetro_backend.app:app
Restart=on-failure
RestartSec=3
UMask=0022

[Install]
WantedBy=multi-user.target
