# ops/ansible/templates/mspmetro-ui.service.j2
[Unit]
Description=MSPMetro UI (Axum + Dioxus SSR)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ mspmetro_deploy_user }}
Group={{ mspmetro_deploy_group }}
WorkingDirectory={{ mspmetro_repo_dir }}
Environment=RUST_LOG=mspmetro_ui=info,tower_http=info
EnvironmentFile=-/etc/default/mspmetro-ui
ExecStart=/usr/local/bin/mspmetro-ui
Restart=on-failure
RestartSec=3
UMask=0022

[Install]
WantedBy=multi-user.target
